# check_gateway.cron
*/5 *   * * *   root    /opt/fffd-check-gateway/check_gateway.sh &>/dev/null
