# check_gateway.cron
*/5 *   * * *   root    /opt/fffd-utils/check_gateway.sh &>/dev/null
