# nodestats.cron
* * * * *       root    /opt/fffd-utils/nodestats.py
