# update-repo.cron
# update repos


# update the fffd-utils repo
0 * * * *	root	/opt/fffd-utils/update-repo.sh -n fffd-utils -d /opt/fffd-utils -c true

# update fffd-dns
#*/10 * * * *	root	/opt/fffd-utils/update-repo.sh -n fffd-dns -d /opt/fffd-dns -c "service bind9 reload"

