# update-repo.cron
# update repos


# update the fffd-utils repo
0 * * * *	root	/opt/fffd-utils/update-repo.sh -n fffd-utils -d /opt/fffd-utils -c true

# update the fffd-peers-backbone
*/15 * * * *	root	/opt/fffd-utils/update-repo.sh -n fffd-peers-backbone -d /etc/fastd/fffd-mesh-vpn/fffd-peers-backbone -c "pkill -HUP fastd"


